# edgenode/supervisord.conf
[supervisord]
nodaemon=true

[program:edgenode]
command=bash start_flume.sh
redirect_stderr=true
stdout_logfile=flume.log
